# deploy/buythistweet.socket
[Unit]
Description=buythistweet gunicorn socket

[Socket]
ListenStream=/run/buythistweet/socket

[Install]
WantedBy=sockets.target